# build/roles/europa/files/xscreensaver/xscreensaver.service
[Unit]
Description=XScreenSaver
[Service]
ExecStart=/usr/bin/xscreensaver
[Install]
WantedBy=default.target