# build/roles/europa/files/openshift/openshift.service
[Unit]
Description=OpenShift Container Application Platform
After=syslog.target
After=network.target

[Service]
WorkingDirectory=/usr/local/openshift/default
EnvironmentFile=-/etc/systemd/system/openshift.service.d/openshift.conf
ExecStart=/usr/local/openshift/default/openshift start
Restart=always

[Install]
WantedBy=multi-user.target