# build/roles/europa/templates/xscreensaver/xscreensaver.service.j2
[Unit]
Description=XScreenSaver
[Service]
ExecStart=/usr/bin/sudo -u {{user}} /usr/bin/xscreensaver
[Install]
WantedBy=default.target